# docker/scitq-server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:scitq-queue]
environment=QUEUE_PROCESS=1
command=scitq-queue
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:scitq-main]
command=pyuwsgi --http :5000 --http-websockets --master -p %(ENV_SCITQ_SERVER_PROCESS)s --enable-threads -w scitq.server:app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true